# deploy/edge/quadlets/cloudflared.container
[Unit]
Description=Cloudflare Tunnel (public access to the demo app)
# Origin is Traefik; start after it so the tunnel doesn't come up pointing at nothing.
After=traefik.service
Wants=traefik.service
# Power-loss boot race: on cold start cloudflared can launch before upstream DNS/internet
# is reachable, fail edge discovery, and exhaust the default 5-in-10s restart burst —
# parking the unit in `failed` so it never retries even once the network is up.
# Disable the rate limiter so Restart=always keeps trying indefinitely until the edge resolves.
StartLimitIntervalSec=0

[Container]
ContainerName=cloudflared
# TODO: pin to a specific release tag to match repo convention (cf. traefik:v3.3, mosquitto:2.0.21).
Image=docker.io/cloudflare/cloudflared:latest
# Must share trakrf.network so the dashboard/Terraform ingress can target https://traefik:443 by name.
Network=trakrf.network
# Token only — kept out of the shared .env so DB/JWT secrets aren't exposed to this sidecar.
# Managed by trakrf/infra Terraform (TRA-957); drop the real value in /srv/trakrf/secrets/cloudflared.env.
EnvironmentFile=/srv/trakrf/secrets/cloudflared.env
# Outbound-only: no PublishPort (the whole point — no inbound, NAT/double-NAT agnostic).
# cloudflared reads the tunnel token from $TUNNEL_TOKEN. Public hostname -> service ingress
# (app.demo.trakrf.id -> https://traefik:443) is configured in Cloudflare, not here.
Exec=tunnel --no-autoupdate run

[Service]
Restart=always
# Wait between attempts so retries don't burn instantly; with StartLimitIntervalSec=0
# this loops patiently through a slow network-up after a power loss.
RestartSec=10

[Install]
WantedBy=default.target
