# deploy/edge/quadlets/mosquitto.container
[Unit]
Description=Mosquitto broker (demo box)

[Container]
ContainerName=mosquitto
Image=docker.io/library/eclipse-mosquitto:2.0.21
Network=trakrf.network
Volume=/srv/trakrf/config/mosquitto/mosquitto.conf:/mosquitto/config/mosquitto.conf:ro,Z
Volume=/srv/trakrf/secrets/mosquitto/passwd:/mosquitto/config/passwd:ro,Z
PublishPort=1883:1883

[Service]
Restart=always

[Install]
WantedBy=default.target
