# deploy/edge/quadlets/traefik.container
[Unit]
Description=Traefik TLS edge
After=backend.service
Requires=backend.service

[Container]
ContainerName=traefik
Image=docker.io/library/traefik:v3.3
Network=trakrf.network
Volume=/srv/trakrf/config/traefik/traefik.yaml:/etc/traefik/traefik.yaml:ro,Z
Volume=/srv/trakrf/config/traefik/dynamic.yaml:/etc/traefik/dynamic.yaml:ro,Z
Volume=/srv/trakrf/secrets/traefik/certs:/certs:ro,Z
PublishPort=443:443

[Service]
Restart=always

[Install]
WantedBy=default.target
