# deploy/edge/systemd/rootlessport-watchdog.service
[Unit]
Description=TrakRF rootlessport forward watchdog (mosquitto + timescaledb)

[Service]
Type=oneshot
ExecStart=/srv/trakrf/scripts/rootlessport-watchdog.sh
